# rootfs/etc/services.d/crond/finish
#!/usr/bin/with-contenv sh

exec /bin/true
